# install/pi-sensor.service.tmpl
[Unit]
Description=Pi-Sensor
After=network.target

[Service]
ExecStart=/home/pi/pi-sensor
WorkingDirectory=/home/pi/
StandardOutput=inherit
StandardError=inherit
Restart=always
User=pi
Environment=CLOUDMQTT_URL={{.CLOUDMQTT_URL}}
Environment=TOPIC={{.TOPIC}}
Environment=SENSOR_SOURCE={{.SENSOR_SOURCE}}

[Install]
WantedBy=multi-user.target
